# ralph-notion.service
[Unit]
Description=Ralph Notion Loop - Autonomous Claude Code Task Runner
After=network.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__
ExecStart=__RALPH_HOME__/ralph_loop.sh --notion
Restart=on-failure
RestartSec=30
Environment=HOME=__USER_HOME__

# Logging
StandardOutput=append:/var/log/ralph-notion.log
StandardError=append:/var/log/ralph-notion.log

[Install]
WantedBy=multi-user.target
